# agent-os.service
[Unit]
Description=Agent OS - governed agent runtime + web console
Documentation=https://github.com/InstaWP/agent-os
After=network.target
Wants=network.target

[Service]
Type=simple
User=vikas
Group=vikas
WorkingDirectory=/home/vikas/tools/agent-os
Environment=NODE_ENV=production
Environment=PORT=3010
# node is installed via nvm on this host; point at the absolute binary.
# If node is upgraded, update this path (and PATH below).
# ~/.local/bin holds the `claude` CLI that claude-code agents launch.
Environment=PATH=/home/vikas/.local/bin:/home/vikas/.nvm/versions/node/v22.22.0/bin:/usr/bin:/usr/local/bin
ExecStart=/home/vikas/.nvm/versions/node/v22.22.0/bin/node dist/cli.js serve
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=10
SyslogIdentifier=agent-os
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/vikas/tools/agent-os/data
ReadWritePaths=/tmp
# claude-code agents run the `claude` CLI, which writes session state/transcripts here.
ReadWritePaths=/home/vikas/.claude

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
