# agent-os.service
[Unit]
Description=Agent OS - governed agent runtime + web console
Documentation=https://github.com/InstaWP/agent-os
After=network.target
Wants=network.target

[Service]
Type=simple
User=vikas
Group=vikas
WorkingDirectory=/home/vikas/tools/agent-os
Environment=NODE_ENV=production
Environment=PORT=3010
# node is installed via nvm on this host; point at the absolute binary.
# If node is upgraded, update this path (and PATH below).
# ~/.local/bin holds the `claude` CLI that claude-code agents launch.
Environment=PATH=/home/vikas/.local/bin:/home/vikas/.nvm/versions/node/v22.22.0/bin:/usr/bin:/usr/local/bin
ExecStart=/home/vikas/.nvm/versions/node/v22.22.0/bin/node dist/cli.js serve
ExecReload=/bin/kill -HUP $MAINPID
# KillMode=process (NOT mixed/control-group): on stop/restart, systemd signals ONLY the main node
# process and leaves the rest of the cgroup alone. Agent sessions run in a tmux server that daemonises
# out of node's process tree — but a double-fork escapes the process tree, NOT the cgroup, so a
# cgroup-wide kill (mixed/control-group) SIGKILLs the tmux server on every restart and every live agent
# session dies (they resurface as `crashed`). This is the Linux-only session-loss bug macOS/launchd
# never had (launchd has no cgroups). With `process`, the tmux server survives a restart and the fresh
# node process re-adopts the still-live sessions via the persistent <home>/tmux.sock. Trade-off (intended,
# and matches macOS): a full `stop` that stays down leaves sessions running orphaned until the service
# returns. ttyd/Slack/etc. are still reaped explicitly by the SIGINT handler (registry.stopAll).
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=10
SyslogIdentifier=agent-os
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
# PrivateTmp MUST be false here — it is incompatible with KillMode=process (above). Agent sessions run
# in a tmux server that we deliberately keep alive across a restart (KillMode=process). PrivateTmp gives
# each service *invocation* its own throwaway /tmp; on restart systemd tears the old invocation's /tmp
# down, but the surviving tmux server is still pinned to that now-dead mount namespace — so the `claude`
# CLI's `mkdir /tmp/claude-<uid>` fails with ENOENT and every session dies on the first restart-after-fix.
# With PrivateTmp=false the service shares the host's stable /tmp, which persists across restarts, so a
# surviving session keeps a valid /tmp. (macOS/launchd has no private-tmp/cgroup story — this is the
# Linux-only pairing needed to match its "sessions survive a restart" behaviour.)
PrivateTmp=false
ProtectSystem=strict
# ProtectHome=read-only would be tighter, BUT the OS pre-accepts Claude Code's one-time folder-trust
# dialog per agent folder by seeding ~/.claude.json (a file in the HOME ROOT, written via atomic
# temp+rename → needs the home DIRECTORY writable, not just ~/.claude). Under ProtectHome=read-only with
# only ~/.claude carved out, that write fails silently and every INTERACTIVE session then hangs forever
# on "Do you trust the files in this folder?". So we make the service user's home writable and re-lock
# the sensitive bits below (ProtectSystem=strict still shields /usr, /etc, /boot). If you prefer
# ProtectHome=read-only, you must additionally make the home root writable enough for ~/.claude.json.
ReadWritePaths=/home/vikas
ReadWritePaths=/tmp
# Re-lock the real persistence vectors a compromised agent could abuse (SSH keys, shell rc).
ReadOnlyPaths=/home/vikas/.ssh
# NOTE (fresh box): the ~/.claude dir must exist AND hold a completed login (`claude` → /login, or
# ANTHROPIC_API_KEY) — otherwise every spawned session launches claude, hits "Not logged in", and dies
# instantly (symptom: "sessions won't spawn"). Any ReadWritePaths/ReadOnlyPaths path must already exist
# or systemd fails start with status=226/NAMESPACE.

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
